# deploy/portfolio.service
# systemd unit template for the engineering-portfolio static site.
# __USER__ / __WORKDIR__ are substituted at install time by vps-setup/bootstrap.sh,
# or by hand:  sed -e "s|__USER__|$USER|g" -e "s|__WORKDIR__|$PWD|g" deploy/portfolio.service \
#              | sudo tee /etc/systemd/system/portfolio.service
[Unit]
Description=engineering-portfolio live site — 0.0.0.0:8901 (open 8901 in the provider firewall)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/python3 -m http.server 8901 --bind 0.0.0.0
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
